# docker-scripts/etc/services.d/ollama/run
#!/usr/bin/with-contenv sh
echo "S6: Starting Ollama..."
exec /usr/bin/ollama serve
